# ansible/roles/virtualbmc-daemon/templates/vbmcd.service.j2
[Unit]
Description=Virtual BMC daemon

[Service]
Type=simple
Restart=on-failure
ExecStart="{{ vbmcd_virtualenv_path }}/bin/vbmcd" --foreground
