# ansible/templates/vbmcd.service.j2
[Unit]
Description=Virtual BMC daemon

[Service]
Type=simple
Restart=on-failure
ExecStart="{{ virtualenv_path }}/bin/vbmcd" --foreground
